# crontab
*  *  *  *  *    root    switch_the_relay 3 30
